# deploy/flask-api.service
[Unit]
Description=MP-DH415 Flask API (frontend/api/app.py)
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__REPO_DIR__/frontend/api
Environment="PATH=__VENV_DIR__/bin"
EnvironmentFile=__REPO_DIR__/.env
ExecStart=__VENV_DIR__/bin/gunicorn --bind 0.0.0.0:5001 --workers 2 --timeout 120 app:app
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
